# deploy/oci/stockman-cron.service
[Unit]
Description=Stockman scheduled automation
After=stockman.service
Requires=stockman.service

[Service]
Type=oneshot
User=stockman
WorkingDirectory=/opt/stockman/current
EnvironmentFile=/etc/stockman/stockman.env
ExecStart=/opt/stockman/current/scripts/oci-cron.sh
